# Docker/ovs-docker/supervisord.conf
[supervisord]
nodaemon=true

[program:openvswitch-switch]
command=service openvswitch-switch start

[program:shell]
command=/bin/bash
